# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Backend: migrate then serve ───────────────────────────────────────────────
[program:backend]
command=/bin/sh -c "cd /app/backend && alembic upgrade head && exec uvicorn app.main:app --host 127.0.0.1 --port 8000"
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Frontend: Next.js standalone server ───────────────────────────────────────
[program:frontend]
command=node /app/frontend/server.js
environment=PORT=3000,HOSTNAME=0.0.0.0,NODE_ENV=production
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
